# cron/crontab
0 4 * * * /root/pg_dump.sh >> /var/log/pg_dump.log 2>&1